# env-bismark/wdr3600/files/etc/uci-defaults/55-dnsmasq
#!/bin/sh

if [ -z `uci -q get dhcp.@dnsmasq[0].interface_name` ]; then
	uci -q set dhcp.@dnsmasq[0].interface_name=myrouter.projectbismark.net,br-lan
	uci commit
fi
