# env-bismark/wndr3x00/files/etc/uci-defaults/10-disable-nodogsplash
#!/bin/sh

/etc/init.d/nodogsplash disable
