# env-bismark/wndr3x00/files/etc/uci-defaults/65-freifunk
#!/bin/sh

[ -n "$(uci -q show freifunk.community)" ] || {
  uci set freifunk.community=community
  uci set freifunk.community.name=bismark
  uci commit
}
