# env-bismark/wndr3x00/files/etc/uci-defaults/network-setlanip
#!/bin/sh

if [ "$(uci -q get network.lan.ipaddr)" != '192.168.142.1' ]; then
	uci -q set network.lan.ipaddr=192.168.142.1
	uci commit
fi
